# platform/docker/supervisord/supervisord.conf
; Base supervisord configuration for the unified platform image.
; The entrypoint treats this file as a runtime template and may update it before
; starting supervisord:
; - replaces the backend DATABASE_URL placeholder with the resolved connection URL
; - appends the postgres program when the bundled internal database is used
;
; ngrok tunneling (ARCHESTRA_NGROK_AUTH_TOKEN) is handled in-process by the
; backend via the ngrok agent SDK, so no supervisord program is added for it.

; Important notes about supervisord environment variable handling:
; - Subprocesses automatically inherit all environment variables from the shell
;   that starts supervisord:
;   https://supervisord.org/subprocess.html#subprocess-environment
; - You only need to explicitly set environment variables in program configs when:
;   1. Using placeholder substitution, e.g. DATABASE_URL="placeholder" replaced
;      by the entrypoint at runtime
;   2. Mapping ARCHESTRA_* variables to different names like NEXT_PUBLIC_*
;   3. Setting hardcoded values not already present in the environment
; - Environment variables declared in the Docker image are otherwise inherited
;   automatically by all supervisord-managed programs:
;   https://supervisord.org/configuration.html#supervisord-section-values

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:backend]
directory=/app/backend
# Maintenance mode must be able to boot when the external database is down for
# planned work, so the Docker entrypoint skips migrations while the message is set.
command=/bin/sh -c "if [ -n \"$ARCHESTRA_MAINTENANCE_MODE_MESSAGE\" ]; then ./scripts/start-production.sh; else sleep 5 && ./node_modules/.bin/drizzle-kit migrate && ./scripts/start-production.sh; fi"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10
# DATABASE_URL is set to a placeholder here and replaced via sed in docker-entrypoint.sh
# All other environment variables (NODE_ENV, ARCHESTRA_*, etc.) are inherited automatically
environment=DATABASE_URL="placeholder"

[program:frontend]
directory=/app/frontend
command=/bin/sh -c "sleep 8 && node server.js"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
# HOSTNAME="0.0.0.0" is required for Next.js standalone server to bind to all interfaces.
# Without this, Next.js binds only to the container's IP, breaking kubectl port-forward.
# See: https://nextjs.org/docs/app/api-reference/config/next-config-js/output#automatically-copying-traced-files
# NEXT_PUBLIC_* variables are mapped from ARCHESTRA_* ENV variables
# All other environment variables (NODE_ENV, etc.) are inherited automatically
environment=HOSTNAME="0.0.0.0",NEXT_PUBLIC_ARCHESTRA_API_BASE_URL="%(ENV_ARCHESTRA_API_BASE_URL)s",NEXT_PUBLIC_ARCHESTRA_SENTRY_FRONTEND_DSN="%(ENV_ARCHESTRA_SENTRY_FRONTEND_DSN)s",NEXT_PUBLIC_ARCHESTRA_SENTRY_ENVIRONMENT="%(ENV_ARCHESTRA_SENTRY_ENVIRONMENT)s",NEXT_PUBLIC_ARCHESTRA_ENTERPRISE_LICENSE_ACTIVATED="%(ENV_ARCHESTRA_ENTERPRISE_LICENSE_ACTIVATED)s",NEXT_PUBLIC_ARCHESTRA_ENTERPRISE_LICENSE_FULL_WHITE_LABELING="%(ENV_ARCHESTRA_ENTERPRISE_LICENSE_FULL_WHITE_LABELING)s"

[program:startup-banner]
command=/bin/sh -c "sleep 12 && /app/docker-banner.sh"
autostart=true
autorestart=false
startsecs=0
priority=999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
